# service/com.ane-lm.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ane-lm.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/ane-lm-server</string>
    <string>--model</string>
    <string>/path/to/your/Qwen3.5-0.8B</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8080</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/tmp/ane-lm-server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/ane-lm-server.err.log</string>
</dict>
</plist>
